# data/Dockerfiles/mailman-core/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:mailman]
command=/usr/local/bin/master --force
user=mailman
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:mailman-log]
command=/usr/bin/tail -f /opt/mailman/var/logs/mailman.log
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2
