# data/Dockerfiles/dovecot/supervisord.conf
[supervisord]
nodaemon=true

[program:syslog-ng]
command=/usr/sbin/syslog-ng --foreground --no-caps
redirect_stderr=true
autostart=true
stdout_syslog=true

[program:dovecot]
command=/usr/sbin/dovecot -F
autorestart=true

[program:logfiles]
command=/usr/bin/tail -f /var/log/mail.log /var/log/syslog
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:cron]
command=/usr/sbin/cron -f
autorestart=true
