# data/Dockerfiles/mailman-web/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:uwsgi]
command=uwsgi --ini /opt/mm_web/uwsgi.ini
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2

[program:statics]
command=/usr/bin/python -u -m SimpleHTTPServer 9192
directory=/opt/mm_web-data/
user=mailman
autorestart=true
priority=3

[program:mailman-log]
command=/usr/bin/tail -f /opt/mm_web/uwsgi.log /opt/mm_web/django.log
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2
