# data/Dockerfiles/policyd/files/etc/supervisor/conf.d/cbpolicyd.conf
[program:cbpolicyd]
command=/usr/sbin/cbpolicyd --fg --config=/etc/cluebringer/cluebringer.conf
autostart=true
autorestart=true
startretries=3

# Redirect output so we can see it using "docker logs"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

