# data/Dockerfiles/sogo/supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
autostart=true
autorestart=true
redirect_stderr=true

[program:sogo]
command=/usr/sbin/sogod
user=sogo
autorestart=true

[program:sogo-syslog]
command=/usr/bin/tail -f /var/log/syslog -f /var/log/sogo/sogo.log
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:sogo-cron]
command=/sogo-cron.sh
user=sogo
autorestart=true
