# data/Dockerfiles/sogo/supervisord.conf
[supervisord]
nodaemon=true

[program:syslog-ng]
command=/usr/sbin/syslog-ng --foreground  --no-caps
redirect_stderr=true
autostart=true
stdout_syslog=true

[program:sogo]
command=/usr/sbin/sogod
user=sogo
autorestart=true
priority=20

[program:reconf-domains]
command=/reconf-domains.sh
autorestart=true
priority=10

[program:sogo-syslog]
command=/usr/bin/tail -f /var/log/syslog -f /var/log/sogo/sogo.log
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:cron]
command=/usr/sbin/cron -f
autorestart=true

[program:sogo-webres]
command=/usr/bin/python -u -m SimpleHTTPServer 9192
directory=/usr/lib/GNUstep/SOGo/
user=sogo
autorestart=true

[inet_http_server]
port=9191

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://localhost:9191
