# data/Dockerfiles/postfix/supervisord.conf
[supervisord]
nodaemon=true

[program:syslog-ng]
command=/usr/sbin/syslog-ng --foreground  --no-caps
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true

[program:postfix]
command=/opt/postfix.sh
autorestart=true

[unix_http_server]
file=/var/tmp/supervisord.sock  
chmod=0770  
chown=nobody:nogroup

[supervisorctl]
serverurl=unix:///var/tmp/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
