# data/Dockerfiles/postfix/supervisord.conf
[supervisord]
nodaemon=true

[program:syslog-ng]
command=/usr/sbin/syslog-ng --foreground  --no-caps
redirect_stderr=true
autostart=true
stdout_syslog=true

[program:postfix]
command=/opt/postfix.sh
autorestart=true

[program:postfix-maillog]
command=/bin/tail -f /var/log/zeyple.log /var/log/combined.log
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[unix_http_server]
file=/var/tmp/supervisord.sock  
chmod=0770  
chown=nobody:nogroup

[supervisorctl]
serverurl=unix:///var/tmp/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
